# scripts/.local/bin/protonmail-bridge
#!/sbin/openrc-run
# Copyright 2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2
# by ivan

name="protonmail-bridge daemon"
description=""
command=/usr/bin/proton-mail-bridge
#command_args="${protonmail-bridge_args}"
command_args="--no-window --noninteractive"
command_user="protonmail:protonmail"
# protonmail-bridge does not daemonize itself
command_background=true
pidfile=/run/protonmail-bridge.pid

depend() {
	need net
	use dns
}
